# build/kafka-image/supervisord.conf
[program:zookeeper]
stdout_logfile = /dev/stdout
stderr_logfile = /dev/stderr
command = bin/zookeeper-server-start.sh config/zookeeper.properties

[program:kafka-config]
autorestart = false
startretries = 0
stdout_logfile = /dev/stdout
stderr_logfile = /dev/stderr
command = bin/kafka-docker-config.sh

[program:kafka]
depends_on = zookeeper, kafka-config
stdout_logfile = /dev/stdout
stderr_logfile = /dev/stderr
command = bin/kafka-server-start.sh config/server.properties
